# config/etc/supervisor.d/conf-available/amavisd.conf
;
;This file is managed by iRedMail Team <support@iredmail.org> with Ansible,
;please do __NOT__ modify it manually.
;

[program:amavisd]
command=/usr/sbin/amavisd-new -u amavis -g amavis -c /etc/amavis/conf.d/50-user foreground
priority=30
startsecs=0
autostart=true
autorestart=true
stdout_syslog=true
stderr_syslog=true
