# config/etc/supervisor.d/conf-available/phpfpm.conf
;
;This file is managed by iRedMail Team <support@iredmail.org> with Ansible,
;please do __NOT__ modify it manually.
;

[program:phpfpm]
command=php-fpm7 -c /etc/php7 -F
priority=30
startsecs=0
autostart=true
autorestart=true
stdout_syslog=true
stderr_syslog=true
