# config/var/spool/cron/crontabs/root
#Ansible: Scan spam/ham reported by end users.
*/10 * * * * root /bin/bash /usr/local/bin/scan_reported_mails >/dev/null
#Ansible: iRedAPD: Clean up expired tracking records hourly.
1 * * * * root /usr/bin/python2 /opt/iredapd/tools/cleanup_db.py >/dev/null
#Ansible: iRedAPD: Convert SPF DNS record of specified domain names to IP
*/30 * * * * root /usr/bin/python2 /opt/iredapd/tools/spf_to_greylist_whitelists.py >/dev/null
#Ansible: Roundcube: cleanup Roundcube SQL database
2 2 * * * root /usr/bin/php /opt/www/roundcubemail/bin/cleandb.sh >/dev/null
#Ansible: Roundcube: cleanup Roundcube temporary files under 'temp/' directory
3 2 * * * root /usr/bin/php /opt/www/roundcubemail/bin/gc.sh >/dev/null
