# config/etc/supervisor/conf-available/cron.conf
;
;This file is managed by iRedMail Team <support@iredmail.org> with Ansible,
;please do __NOT__ modify it manually.
;

[program:cron]
command=/usr/sbin/cron -f
priority=5
startsecs=0
autostart=true
autorestart=true
stdout_syslog=true
stderr_syslog=true
