# config/etc/supervisor/conf-available/dovecot.conf
;
;This file is managed by iRedMail Team <support@iredmail.org> with Ansible,
;please do __NOT__ modify it manually.
;

[program:dovecot]
command=/usr/sbin/dovecot -F -c /etc/dovecot/dovecot.conf
priority=20
startsecs=0
autostart=true
autorestart=true
stdout_syslog=true
stderr_syslog=true
