# config/etc/supervisor/conf-available/iredadmin.conf
;
;This file is managed by iRedMail Team <support@iredmail.org> with Ansible,
;please do __NOT__ modify it manually.
;

[program:iredadmin]
command=/usr/bin/uwsgi --ini /opt/www/iredadmin/rc_scripts/uwsgi/debian.ini
priority=20
startsecs=0
autostart=true
autorestart=true
stdout_syslog=true
stderr_syslog=true
