# config/etc/supervisor/conf-available/sogo.conf
;
;This file is managed by iRedMail Team <support@iredmail.org> with Ansible,
;please do __NOT__ modify it manually.
;

[program:sogo]
command=/usr/sbin/sogod -WOUseWatchDog NO -WOLogFile - -WONoDetach YES -WOWorkersCount 3
priority=50
startsecs=0
autostart=true
autorestart=true
stdout_syslog=true
stderr_syslog=true
